# scripts/bitharbor.service
[Unit]
Description=BitHarbor FastAPI Service
After=network-online.target
Wants=network-online.target

[Service]
User=bitharbor
Group=bitharbor
WorkingDirectory=/home/ethan/documents/bitharbor
Environment="BIT_HARBOR_CONFIG=/etc/bitharbor/config.yaml"
ExecStart=/opt/bitharbor/venv/bin/uvicorn bitharbor.main:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5
KillMode=control-group
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

